# docker/supervisord.conf
# NovaMind Supervisord - non-root main config
# No [unix_http_server] (would need a socket in root-owned /run).
# pidfile under /tmp; container runs as appuser (see Dockerfile USER).

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/supervisor

[program:nginx]
command=nginx -c /etc/nginx/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx-access.log
stderr_logfile=/var/log/supervisor/nginx-error.log

[program:backend]
command=python main.py --config docker --host 127.0.0.1 --port 8100
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/backend-access.log
stderr_logfile=/var/log/supervisor/backend-error.log
startsecs=5